# systemd-units/secrets-decrypt.service
[Unit]
Description=Decrypt secrets into tmpfs
Before=network.target
Wants=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/decrypt-secrets.sh
RemainAfterExit=yes
# Запускаем от root для доступа к зашифрованным секретам
User=root
Group=root

# Ограничения безопасности
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
RestrictNamespaces=yes

[Install]
WantedBy=multi-user.target
